# docker/supervisor.conf
[supervisord]
nodaemon=false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm82 -F 
autostart=true
autorestart=true
